# unix_services/radix-display.service
[Unit]
Description=Radix IoT Display Service
After=network.target

[Service]
Type=simple
User=bigwiz
WorkingDirectory=/home/bigwiz/projects/radixiot/utils
ExecStart=/home/bigwiz/projects/radixiot/utils/venv/bin/python display.py
Restart=always
RestartSec=5
Environment="PATH=/home/bigwiz/projects/radixiot/utils/venv/bin:/usr/bin"
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=radix-display

[Install]
WantedBy=multi-user.target
